# sh/initd-header.sh
#!/bin/sh

### BEGIN INIT INFO
# Provides:          asmodee-werewolf
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: asmodee-werewolf
# Description:       Werewolves of Miller's Hollow
### END INIT INFO
